# voicebot.service
[Unit]
Description=Voicebot SAV Wouippleul - AI Server (Distributed Architecture)
Documentation=https://github.com/votre-org/PY_SAV
After=network.target

[Service]
Type=simple

# User/Group (créer un utilisateur dédié pour la sécurité)
User=voicebot
Group=voicebot

# Répertoires
WorkingDirectory=/opt/PY_SAV
Environment="PATH=/opt/PY_SAV/venv/bin:/usr/local/bin:/usr/bin:/bin"

# Variables d'environnement (ou utiliser EnvironmentFile)
# EnvironmentFile=/opt/PY_SAV/.env

# Commande de démarrage
ExecStart=/opt/PY_SAV/venv/bin/python /opt/PY_SAV/server.py

# Restart automatique
Restart=always
RestartSec=10

# Limites de ressources (optionnel)
LimitNOFILE=65536
LimitNPROC=4096

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=voicebot

# Sécurité (optionnel - renforcement)
# PrivateTmp=true
# NoNewPrivileges=true
# ProtectSystem=strict
# ProtectHome=true
# ReadWritePaths=/opt/PY_SAV/logs /opt/PY_SAV/assets/cache

[Install]
WantedBy=multi-user.target
